# demo/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:ballcone]
command=/usr/local/bin/ballcone -m '/var/lib/ballcone' -wh '0.0.0.0'
