# demo/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:balcone]
command=/usr/local/bin/balcone -m '/var/lib/balcone' -wh '0.0.0.0'
